# bot/daha-bot.service
[Unit]
Description=Daha Telegram Bot
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=daha
WorkingDirectory=/opt/daha/bot
Environment=PATH=/opt/daha/bot/venv/bin
ExecStart=/opt/daha/bot/venv/bin/python main.py
Restart=always
RestartSec=10

# Environment variables
Environment=BOT_TOKEN=your_telegram_bot_token_here
Environment=BOT_DB_URL=postgresql://daha_user:SWP2025@localhost/daha_db

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/daha/bot

[Install]
WantedBy=multi-user.target 